# deploy/systemd/opendray.service
# opendray v2 systemd unit
#
# Standard Linux deploy: one binary at /usr/local/bin/opendray,
# one config at /etc/opendray/config.toml, secrets injected via
# an EnvironmentFile so they never sit in the config TOML.
#
# Install:
#   sudo cp deploy/systemd/opendray.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now opendray
#   sudo systemctl status opendray
#
# Inspect:
#   journalctl -u opendray -f --no-pager
#
# This unit assumes:
#   * /usr/local/bin/opendray exists (binary or symlink).
#   * Service user `opendray` exists with $HOME at /var/lib/opendray.
#   * /etc/opendray/config.toml is mode 600, owned by opendray:opendray.
#   * /etc/opendray/env.d/secrets is mode 600, owned by root or
#     opendray, and contains at minimum:
#         OPENDRAY_ADMIN_PASSWORD=...
#         OPENDRAY_BACKUP_KEY=...   (only if [backup].enabled = true)
#
# See deploy/README.md for the full bootstrap walkthrough.

[Unit]
Description=opendray v2 — multiplexer + integration gateway for AI agent CLIs
Documentation=https://github.com/Opendray/opendray/blob/main/docs/operator-guide.md
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=opendray
Group=opendray
WorkingDirectory=/var/lib/opendray

# Secrets land here; this file is owned root:opendray, mode 0640.
EnvironmentFile=/etc/opendray/env.d/secrets

# Run pending migrations before starting the server. Idempotent;
# re-running is a no-op when the schema is current. Failure here
# blocks the main ExecStart, which is what we want.
ExecStartPre=/usr/local/bin/opendray migrate -config /etc/opendray/config.toml

# Bind to 127.0.0.1 by default (per config.example.toml). Front
# with a reverse proxy or Cloudflare Tunnel for TLS.
ExecStart=/usr/local/bin/opendray serve -config /etc/opendray/config.toml

# Restart policy — opendray returns non-zero on fatal errors, never
# on graceful shutdown. on-failure restarts only when something
# actually broke; SIGTERM/SIGINT exits cleanly with status 0.
#
# RestartForceExitStatus adds SIGHUP/SIGPIPE to the restart set.
# systemd treats SIGHUP/SIGINT/SIGTERM/SIGPIPE as "clean exit" by
# default — opendray catches SIGTERM/SIGINT for graceful shutdown,
# but HUP/PIPE arriving unhandled means something killed us we
# didn't expect (rogue cleanup script, parent shell on logout,
# misconfigured tool). Treat those as failures and restart.
Restart=on-failure
RestartSec=5s
RestartForceExitStatus=SIGHUP SIGPIPE
StartLimitIntervalSec=60
StartLimitBurst=5

# Resource ceilings — the default ulimit (1024 fds) is easy to
# exhaust with concurrent PTY sessions + DB pool + WS clients.
LimitNOFILE=65536
TasksMax=4096

# Graceful shutdown window — opendray's internal supervisor budgets
# 15s for HTTP drain + subsystem stop. Give systemd 20s before
# escalating to SIGKILL.
TimeoutStopSec=20s
KillSignal=SIGTERM
SendSIGKILL=yes

# Sandboxing — adapt for your environment. The defaults below assume
# OpenDray talks only to: /etc/opendray (read-only config) and
# /var/lib/opendray (writable state). Loosen if you point [vault.root]
# or [backup].local_dir elsewhere.
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelLogs=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
# NOTE: MemoryDenyWriteExecute is intentionally NOT enabled.
#
# It enforces W^X on the process AND all its children via seccomp.
# Codex and Gemini are V8/Node.js-based CLIs whose JIT needs to flip
# JIT pages from RW→RX via mprotect, which the W^X filter blocks
# (SIGSYS). The child dies within ~50 ms of pty.Start and the user
# sees a "[session ended — read-only buffer]" banner with no other
# diagnostic.
#
# (Claude survives because anthropic-ai/claude-code falls back to an
# interpreter-only path when it detects the constraint; codex and
# gemini don't have that fallback.)
#
# If you want this hardening back, you'd have to spawn child CLIs in
# a different unit/cgroup so the seccomp filter doesn't inherit.
SystemCallArchitectures=native
ReadWritePaths=/var/lib/opendray

# Capability surface — opendray needs none for the default loopback
# bind on a port >1024. If you bind to a privileged port (<1024)
# directly, add `AmbientCapabilities=CAP_NET_BIND_SERVICE` and
# `CapabilityBoundingSet=CAP_NET_BIND_SERVICE`.
CapabilityBoundingSet=

[Install]
WantedBy=multi-user.target
